# beaglebone/copter-data.conf
# Stored in /etc/init/copter-data.conf on the Beaglebone Debian
description "Collecting data from the copter"

start on filesystem
stop on shutdown

script
    exec bash -c '/etc/copter'
end script
